# deploy/com.crawlhub.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crawlhub.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Volumes/KESU-1/Data-Fin/股东报告/scheduler_service.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/KESU-1/Data-Fin/股东报告</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Volumes/KESU-1/Data-Fin/股东报告/data/scheduler_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/KESU-1/Data-Fin/股东报告/data/scheduler_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>PYTHONIOENCODING</key>
        <string>utf-8</string>
    </dict>
</dict>
</plist>
